# systemd/metabridge-api.service
[Unit]
Description=Metabridge API Server
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/projects/metabridge-engine-hub
Environment="PATH=/usr/local/go/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/root/projects/metabridge-engine-hub/bin/api -config config/config.production.yaml
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
